# integrations/samantha/claude-bridge.conf
[program:claude-bridge]
command=/usr/local/lib/hermes-agent/venv/bin/python3 -m uvicorn claude_bridge:app --host 127.0.0.1 --port 8787
directory=/usr/local/lib/hermes-agent
autostart=true
autorestart=true
startsecs=5
startretries=5
environment=HOME="/root",PATH="/usr/local/bin:/usr/bin:/bin"
redirect_stderr=true
stdout_logfile=/tmp/claude-bridge.log
stdout_logfile_maxbytes=1MB
