# integrations/samantha/crond.conf
[program:crond]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stdout_logfile=/tmp/cron.log
stderr_logfile=/tmp/cron.log
