# integrations/samantha/fix-dns.conf
[program:fix-dns]
command=/bin/bash -c "rm -f /etc/resolv.conf && printf 'nameserver 8.8.8.8\nnameserver 8.8.4.4\n' > /etc/resolv.conf && echo DNS fixed"
autostart=true
autorestart=false
startsecs=0
priority=1
stdout_logfile=/tmp/fix-dns.log
stderr_logfile=/tmp/fix-dns.log
