# integrations/samantha/restore-codex-auth.conf
[program:restore-codex-auth]
command=/root/restore-codex-auth.sh
autostart=true
autorestart=false
startsecs=0
priority=2
stdout_logfile=/tmp/restore-codex-auth.log
stderr_logfile=/tmp/restore-codex-auth.log
