# integrations/samantha/sshd.conf
[program:sshd]
command=/bin/bash -c "mkdir -p /run/sshd && exec /usr/sbin/sshd -D"
autostart=true
autorestart=true
startretries=5
stdout_logfile=/tmp/sshd.log
stderr_logfile=/tmp/sshd.err
priority=5
